# nrf_flasher/rootfs/etc/services.d/flasher/finish
#!/command/with-contenv bash
set -euo pipefail
echo "[dfu] service stopping"